# exercises/ex7.conf
[supervisord]
nodaemon=true

[program:find_command]
command=/bin/bash -c "tail -f /dev/null"
user=user1
directory=/app/user1
autostart=true
autorestart=true
