# exercises/ex8.conf
[supervisord]
nodaemon=true

[program:find_command]
command=/bin/bash -c "stress-ng --fd 8 1>/dev/null 2>/dev/null"
user=user1
directory=/app/user1
autostart=true
autorestart=true
